# postgresql_cluster/templates/postgresql.service.j2
.include /usr/lib/systemd/system/postgresql-{{ cluster.postgresql_version }}.service

[Service]
Environment=PGDATA={{ cluster.path }}
